# scrapper/outputs/threat_actors_data/OilRig.json
[
    {
        "Domain": "Enterprise",
        "ID": "T1087.001",
        "Technique": "Account Discovery: Local Account",
        "Use": "OilRig has run net user , net user /domain , net group \"domain admins\" /domain , and net group \"Exchange Trusted Subsystem\" /domain to get account listings on a victim. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1087.002",
        "Technique": "Account Discovery: Domain Account",
        "Use": "OilRig has run net user , net user /domain , net group \"domain admins\" /domain , and net group \"Exchange Trusted Subsystem\" /domain to get account listings on a victim. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1071.001",
        "Technique": "Application Layer Protocol: Web Protocols",
        "Use": "OilRig has used HTTP for C2. [6] [12] [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1071.004",
        "Technique": "Application Layer Protocol: DNS",
        "Use": "OilRig has used DNS for C2 including the publicly available requestbin.net tunneling service. [6] [12] [13] [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1119",
        "Technique": "Automated Collection",
        "Use": "OilRig has used automated collection. [6]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1110",
        "Technique": "Brute Force",
        "Use": "OilRig has used brute force techniques to obtain credentials. [12]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1059",
        "Technique": "Command and Scripting Interpreter",
        "Use": "OilRig has used various types of scripting for execution. [1] [14] [15] [7] [16]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1071.001",
        "Technique": "PowerShell",
        "Use": "OilRig has used PowerShell scripts for execution, including use of a macro to run a PowerShell command to decode file contents. [1] [17] [9]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1071.003",
        "Technique": "Windows Command Shell",
        "Use": "OilRig has used macros to deliver malware such as QUADAGENT and OopsIE . [1] [14] [15] [7] [16] OilRig has used batch scripts. [1] [14] [15] [7] [16]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1071.005",
        "Technique": "Visual Basic",
        "Use": "OilRig has used VBScript macros for execution on compromised hosts. [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1555",
        "Technique": "Credentials from Password Stores",
        "Use": "OilRig has used credential dumping tools such as LaZagne to steal credentials to accounts logged into the compromised system and to Outlook Web Access. [6] [12] [18] [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1071.003",
        "Technique": "Credentials from Web Browsers",
        "Use": "OilRig has used credential dumping tools such as LaZagne to steal credentials to accounts logged into the compromised system and to Outlook Web Access. [6] [12] [18] [13] OilRig has also used tool named PICKPOCKET to dump passwords from web browsers. [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1071.004",
        "Technique": "Windows Credential Manager",
        "Use": "OilRig has used credential dumping tool named VALUEVAULT to steal credentials from the Windows Credential Manager. [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1140",
        "Technique": "Deobfuscate/Decode Files or Information",
        "Use": "A OilRig macro has run a PowerShell command to decode file contents. OilRig has also used certutil to decode base64-encoded files on victims. [1] [17] [15] [19]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1573.002",
        "Technique": "Encrypted Channel: Asymmetric Cryptography",
        "Use": "OilRig used the Plink utility and other tools to create tunnels to C2 servers. [12]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1048.003",
        "Technique": "Exfiltration Over Alternative Protocol: Exfiltration Over Unencrypted Non-C2 Protocol",
        "Use": "OilRig has exfiltrated data over FTP separately from its primary C2 channel over DNS. [5]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1133",
        "Technique": "External Remote Services",
        "Use": "OilRig uses remote services such as VPN, Citrix, or OWA to persist in an environment. [12]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1008",
        "Technique": "Fallback Channels",
        "Use": "OilRig malware ISMAgent falls back to its DNS tunneling mechanism if it is unable to reach the C2 server over HTTP. [14]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1070.004",
        "Technique": "Indicator Removal: File Deletion",
        "Use": "OilRig has deleted files associated with their payload after execution. [1] [15]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1105",
        "Technique": "Ingress Tool Transfer",
        "Use": "OilRig can download remote files onto victims. [1]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1056.001",
        "Technique": "Input Capture: Keylogging",
        "Use": "OilRig has used keylogging tools called KEYPUNCH and LONGWATCH. [12] [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1036",
        "Technique": "Masquerading",
        "Use": "OilRig has used .doc file extensions to mask malicious executables. [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1046",
        "Technique": "Network Service Discovery",
        "Use": "OilRig has used the publicly available tool SoftPerfect Network Scanner as well as a custom tool called GOLDIRONY to conduct network scanning. [12]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1027.005",
        "Technique": "Obfuscated Files or Information: Indicator Removal from Tools",
        "Use": "OilRig has tested malware samples to determine AV detection and subsequently modified the samples to ensure AV evasion. [2] [16]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1027.013",
        "Technique": "Obfuscated Files or Information: Encrypted/Encoded File",
        "Use": "OilRig has encrypted and encoded data in its malware, including by using base64. [1] [7] [6] [9] [16]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1137.004",
        "Technique": "Office Application Startup: Outlook Home Page",
        "Use": "OilRig has abused the Outlook Home Page feature for persistence. OilRig has also used CVE-2017-11774 to roll back the initial patch designed to protect against Home Page abuse. [20]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1003.001",
        "Technique": "OS Credential Dumping: LSASS Memory",
        "Use": "OilRig has used credential dumping tools such as Mimikatz to steal credentials to accounts logged into the compromised system and to Outlook Web Access. [6] [12] [18] [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1003.004",
        "Technique": "OS Credential Dumping: LSA Secrets",
        "Use": "OilRig has used credential dumping tools such as LaZagne to steal credentials to accounts logged into the compromised system and to Outlook Web Access. [6] [12] [18] [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1003.005",
        "Technique": "OS Credential Dumping: Cached Domain Credentials",
        "Use": "OilRig has used credential dumping tools such as LaZagne to steal credentials to accounts logged into the compromised system and to Outlook Web Access. [6] [12] [18] [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1201",
        "Technique": "Password Policy Discovery",
        "Use": "OilRig has used net.exe in a script with net accounts /domain to find the password policy of a domain. [21]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1120",
        "Technique": "Peripheral Device Discovery",
        "Use": "OilRig has used tools to identify if a mouse is connected to a targeted system. [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1069.001",
        "Technique": "Permission Groups Discovery: Local Groups",
        "Use": "OilRig has used net localgroup administrators to find local administrators on compromised systems. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1069.002",
        "Technique": "Permission Groups Discovery: Domain Groups",
        "Use": "OilRig has used net group /domain , net group \"domain admins\" /domain , and net group \"Exchange Trusted Subsystem\" /domain to find domain group permission settings. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1566.001",
        "Technique": "Phishing: Spearphishing Attachment",
        "Use": "OilRig has sent spearphising emails with malicious attachments to potential victims using compromised and/or spoofed email accounts. [15] [7] [9]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1566.002",
        "Technique": "Phishing: Spearphishing Link",
        "Use": "OilRig has sent spearphising emails with malicious links to potential victims. [15]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1566.003",
        "Technique": "Phishing: Spearphishing via Service",
        "Use": "OilRig has used LinkedIn to send spearphishing links. [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1057",
        "Technique": "Process Discovery",
        "Use": "OilRig has run tasklist on a victim's machine. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1572",
        "Technique": "Protocol Tunneling",
        "Use": "OilRig has used the Plink utility and other tools to create tunnels to C2 servers. [6] [12] [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1012",
        "Technique": "Query Registry",
        "Use": "OilRig has used reg query \"HKEY_CURRENT_USER\\Software\\Microsoft\\Terminal Server Client\\Default\" on a victim to query the Registry. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1021.001",
        "Technique": "Remote Services: Remote Desktop Protocol",
        "Use": "OilRig has used Remote Desktop Protocol for lateral movement. The group has also used tunneling tools to tunnel RDP into the environment. [6] [12] [19]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1021.004",
        "Technique": "Remote Services: SSH",
        "Use": "OilRig has used Putty to access compromised systems. [6]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1053.005",
        "Technique": "Scheduled Task/Job: Scheduled Task",
        "Use": "OilRig has created scheduled tasks that run a VBScript to execute a payload on victim machines. [15] [7] [13] [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1113",
        "Technique": "Screen Capture",
        "Use": "OilRig has a tool called CANDYKING to capture a screenshot of user's desktop. [12]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1505.003",
        "Technique": "Server Software Component: Web Shell",
        "Use": "OilRig has used web shells, often to maintain access to a victim network. [6] [12] [19]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1218.001",
        "Technique": "System Binary Proxy Execution: Compiled HTML File",
        "Use": "OilRig has used a CHM payload to load and execute another malicious file once delivered to a victim. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1082",
        "Technique": "System Information Discovery",
        "Use": "OilRig has run hostname and systeminfo on a victim. [4] [5] [13] [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1016",
        "Technique": "System Network Configuration Discovery",
        "Use": "OilRig has run ipconfig /all on a victim. [4] [5]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1049",
        "Technique": "System Network Connections Discovery",
        "Use": "OilRig has used netstat -an on a victim to get a listing of network connections. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1033",
        "Technique": "System Owner/User Discovery",
        "Use": "OilRig has run whoami on a victim. [4] [5] [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1007",
        "Technique": "System Service Discovery",
        "Use": "OilRig has used sc query on a victim to gather information about services. [4]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1552.001",
        "Technique": "Unsecured Credentials: Credentials In Files",
        "Use": "OilRig has used credential dumping tools such as LaZagne to steal credentials to accounts logged into the compromised system and to Outlook Web Access. [6] [12] [18] [13]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1204.001",
        "Technique": "User Execution: Malicious Link",
        "Use": "OilRig has delivered malicious links to achieve execution on the target system. [15] [7] [9]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1204.002",
        "Technique": "User Execution: Malicious File",
        "Use": "OilRig has delivered macro-enabled documents that required targets to click the \"enable content\" button to execute the payload on the system. [15] [7] [9] [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1078",
        "Technique": "Valid Accounts",
        "Use": "OilRig has used compromised credentials to access other systems on a victim network. [6] [12] [19]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1497.001",
        "Technique": "Virtualization/Sandbox Evasion: System Checks",
        "Use": "OilRig has used macros to verify if a mouse is connected to a compromised machine. [10]"
    },
    {
        "Domain": "Enterprise",
        "ID": "T1047",
        "Technique": "Windows Management Instrumentation",
        "Use": "OilRig has used WMI for execution. [12]"
    },
    {
        "Domain": "ICS",
        "ID": "T0817",
        "Technique": "Drive-by Compromise",
        "Use": "OilRig has been seen utilizing watering hole attacks to collect credentials which could be used to gain access into ICS networks. [22]"
    },
    {
        "Domain": "ICS",
        "ID": "T0853",
        "Technique": "Scripting",
        "Use": "OilRig has embedded a macro within spearphishing attachments that has been made up of both a VBScript and a PowerShell script. [23]"
    },
    {
        "Domain": "ICS",
        "ID": "T0865",
        "Technique": "Spearphishing Attachment",
        "Use": "OilRig used spearphishing emails with malicious Microsoft Excel spreadsheet attachments. [23]"
    },
    {
        "Domain": "ICS",
        "ID": "T0869",
        "Technique": "Standard Application Layer Protocol",
        "Use": "OilRig communicated with its command and control using HTTP requests. [23]"
    },
    {
        "Domain": "ICS",
        "ID": "T0859",
        "Technique": "Valid Accounts",
        "Use": "OilRig utilized stolen credentials to gain access to victim machines. [24]"
    }
]